# com.spa.httpserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spa.httpserver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8765</string>
        <string>--directory</string>
        <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/spa_http.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/spa_http.err</string>
</dict>
</plist>
