# scripts/com.spa.daily_cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spa.daily_cycle</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>spa_core.paper_trading.cycle_runner</string>
        <string>--verbose</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>
    <!--
        BUG FIX (TELEGRAM_AUDIT 2026-06-18):
        Was: StartInterval 1800 (every 30 minutes) — caused Telegram spam
             because _run_cycle_alerts fired send_red_flag / send_gap_alert
             (no dedup) every 30 min whenever red_flags.json was non-empty.
        Now: StartCalendarInterval 08:00 daily — one cycle per day as intended.
        CLAUDE.md has always said "ежедневно 08:00"; the plist was wrong.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/spa_cycle.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/spa_cycle_err.log</string>
</dict>
</plist>
