# scripts/com.spa.cycle_gap_monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    L2 — Cycle Gap Monitor (каждые 5 минут)
    Module: spa_core/paper_trading/cycle_gap_monitor.py
    Функция: детектирует пропущенные/задержанные daily-циклы (порог 26h).
             Telegram-алерт не чаще 1 раза в сутки (dedup через
             data/cycle_gap_state.json). Читает paper_trading_status.json.
    Пишет:   data/cycle_gap_state.json  (dedup state, atomic)
    CLI:     python3 -m spa_core.paper_trading.cycle_gap_monitor -run
    Python: /Users/yuriikulieshov/miniconda3/bin/python3 (AGENT-001 fix, 2026-06-22)
    Установка:
      bash scripts/install_all_agents.sh
    -->

    <key>Label</key>
    <string>com.spa.cycle_gap_monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>spa_core.paper_trading.cycle_gap_monitor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yuriikulieshov</string>
        <key>PATH</key>
        <string>/Users/yuriikulieshov/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/spa_cycle_gap_monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/spa_cycle_gap_monitor_err.log</string>
</dict>
</plist>
