# scripts/com.spa.morning_digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
    MP-1493 (v11.09): Morning digest Telegram — fires at 08:05 UTC,
    5 minutes after the daily cycle (08:00), so the cycle data is ready.

    Install:
        cp scripts/com.spa.morning_digest.plist ~/Library/LaunchAgents/
        launchctl load ~/Library/LaunchAgents/com.spa.morning_digest.plist

    Uninstall:
        launchctl unload ~/Library/LaunchAgents/com.spa.morning_digest.plist
        rm ~/Library/LaunchAgents/com.spa.morning_digest.plist
-->
<dict>
    <key>Label</key>
    <string>com.spa.morning_digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>spa_core.alerts.morning_digest</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude/logs/morning_digest_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude/logs/morning_digest_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yuriikulieshov/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
