# scripts/com.spa.rwa_safety_board.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    RWA Collateral Safety Board (daily, ADVISORY)
    Module: spa_core.strategy_lab.rwa_backstop.safety_board
    Function: measures the whole tokenized-RWA collateral universe from LIVE DeFiLlama /pools data
              with the §SPA-RRB-validated LiquidationNAVEngine, writes a per-asset
              LIQUID/THIN/REDEMPTION_ONLY/UNSAFE verdict + marketing-vs-LiqNAV gap %
              to data/rwa_safety_board.json (atomic).
    ADVISORY / RESEARCH ONLY: never lends, never trades, never touches the go-live track.
    Schedule: daily 05:50 local (after com.spa.refusal 05:45, before the 08:00 UTC daily_cycle) + RunAtLoad.
    Writes:  data/rwa_safety_board.json (atomic write)
             data/rwa_nav_curve.json   (measured-NAV forward record, one point/UTC day)
             data/rwa_floor_curve.json (risk-free-FLOOR forward record from the LIVE
                                        tokenized-T-bill feed, one point/UTC day) — the
                                        forward audit trail of the rate the RWAFloor benchmark
                                        + rwa_sleeve accrue at. Idempotent, fail-CLOSED, atomic.
    Exit:    0 = success, non-0 = internal error.
    CLI:     python3 -m spa_core.strategy_lab.rwa_backstop.safety_board
    Install:
      cp scripts/com.spa.rwa_safety_board.plist ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.spa.rwa_safety_board.plist
      launchctl kickstart -k gui/$(id -u)/com.spa.rwa_safety_board
    -->

    <key>Label</key>
    <string>com.spa.rwa_safety_board</string>

    <!-- launchd CANNOT exec miniconda-python directly (exit 78 EX_CONFIG);
         runs via bash wrapper (scripts/agent_rwa_safety_board.sh -> agent_template.sh). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yuriikulieshov/Documents/SPA_Claude/scripts/agent_rwa_safety_board.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Не запускать повторно если предыдущий ещё не завершился -->
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/tmp/spa_rwa_safety_board.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/spa_rwa_safety_board.launchd.err</string>
</dict>
</plist>
