# scripts/com.spa.fund-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.spa.fund-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>/Users/yuriikulieshov/Documents/SPA_Claude/scripts/fund_api_server.py</string>
        <string>8766</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPA_DATA_DIR</key>
        <string>/Users/yuriikulieshov/Documents/SPA_Claude/data</string>
    </dict>

    <!-- stdout → /tmp/spa_fund_api.log, stderr → /tmp/spa_fund_api_err.log -->
    <key>StandardOutPath</key>
    <string>/tmp/spa_fund_api.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/spa_fund_api_err.log</string>

    <!-- Запускается по требованию, НЕ как постоянный демон.
         Порт 8766 (8765 занят com.spa.httpserver). -->
    <key>KeepAlive</key>
    <false/>

    <!-- Не запускать автоматически при load; запускать вручную через launchctl start -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
