# scripts/com.spa.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    L1 PERSISTENT DAEMON — Cloudflare Tunnel (token mode)
    Token-based auth: run_cloudflared.sh reads CF_TUNNEL_TOKEN_SPA from Keychain.
    Ingress config lives on Cloudflare edge, no local config.yml needed.
    Install:
      cp scripts/com.spa.cloudflared.plist ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.spa.cloudflared.plist
    Verify:
      launchctl list | grep com.spa.cloudflared
      tail -f /tmp/spa_cloudflared.log
    -->

    <key>Label</key>
    <string>com.spa.cloudflared</string>

    <!-- Запуск через wrapper: устойчив к Intel (/usr/local) и
         Apple Silicon (/opt/homebrew) путям установки cloudflared. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yuriikulieshov/Documents/SPA_Claude/scripts/run_cloudflared.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <!-- HOME needed for Keychain access; PATH for brew-installed binary. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yuriikulieshov</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- L1: перезапуск при падении туннеля -->
    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <!-- Защита от crash loop: не рестартовать быстрее 10 сек -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/spa_cloudflared.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/spa_cloudflared_err.log</string>
</dict>
</plist>
