# launchd/com.spa.tournament_engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spa.tournament_engine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>spa_core.tournament.tournament_engine</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <!-- Run daily at 09:00 UTC (launchd uses local wall-clock time;
         adjust Hour if your timezone differs from UTC) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do NOT start at load time — only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logs live in the repo (survive reboot; /tmp is cleared). agent_health
         judges freshness by the freshest of these two streams; this module logs
         via Python logging → stderr, so the _err log is the live heartbeat. -->
    <key>StandardOutPath</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude/logs/tournament_engine.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude/logs/tournament_engine_err.log</string>

    <!-- Retry throttle: wait 60 s before re-launching on crash -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yuriikulieshov/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/yuriikulieshov</string>
        <!-- Telegram tokens are read from macOS Keychain at runtime;
             no secrets stored here. -->
    </dict>

</dict>
</plist>
