# launchd/com.spa.rules_watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spa.rules_watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuriikulieshov/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>spa_core.monitoring.rules_watchdog</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuriikulieshov/Documents/SPA_Claude</string>

    <!-- Run every 300 seconds (5 minutes) -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Start at load time -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive on failure for up to 3 restarts -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/spa_watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/spa_watchdog_err.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yuriikulieshov/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/yuriikulieshov</string>
    </dict>

</dict>
</plist>
